# res/ansible/templates/gateway-compose.service.j2
[Unit]
Description=WORKSPACE-GATEWAY Compose Stack
After=network-online.target
Wants=network-online.target

[Service]
EnvironmentFile=-{{ project_root }}/.env
Type=simple
WorkingDirectory={{ project_root }}/res/docker
Environment=PATH={{ compose_cmd | dirname }}:{{ podman_bin_dir }}:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
Environment=COMPOSE_HTTP_TIMEOUT=10
ExecStartPre=-{{ compose_cmd }} --podman-path {{ podman_path }} -f {{ compose_file }} down
ExecStartPre=-{{ podman_path }} network rm -f docker_default
ExecStartPre=-/bin/sh -c '{{ podman_path }} ps -aq --filter "status=exited" --filter "status=dead" --filter "status=created" | xargs -r {{ podman_path }} rm -f'
ExecStart={{ compose_cmd }} --podman-path {{ podman_path }} -f {{ compose_file }} up --remove-orphans
ExecStop={{ compose_cmd }} --podman-path {{ podman_path }} -f {{ compose_file }} down
Restart=always
RestartSec=15
TimeoutStartSec=300
TimeoutStopSec=120
StandardOutput=journal
StandardError=journal
SyslogIdentifier=gateway-compose

[Install]
WantedBy=default.target